# install/autostart/docker-app1.desktop
[Desktop Entry]
Type=Application
Exec=/usr/local/bin/startall_user.sh
Hidden=false
NoDisplay=false

Name[en_US]=DockerApp
Name=DockerApp
