# install/autostart/dockerterminal.desktop
[Desktop Entry]
Name=DockerAppTerminal
Comment=Use the docker app command line
Keywords=shell;prompt;command;commandline;
Exec=/bin/lxterminal -T PVZDclient-Terminal --geometry=150x50 -e /usr/local/bin/dockerterminal.sh
Icon=utilities-terminal
Type=Application

Hidden=false
NoDisplay=false
Name[en_US]=dockerterminal.desktop
